# com.opera.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opera.backend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sandeeplingam/VibeCoding/Opera/.venv/bin/uvicorn</string>
        <string>opera.backend.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/sandeeplingam/VibeCoding/Opera</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/sandeeplingam/VibeCoding/Opera/logs/backend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sandeeplingam/VibeCoding/Opera/logs/backend.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
