# com.opera.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opera.frontend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sandeeplingam/VibeCoding/Opera/start-frontend.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/sandeeplingam/VibeCoding/Opera/opera-frontend</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/sandeeplingam/VibeCoding/Opera/logs/frontend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/sandeeplingam/VibeCoding/Opera/logs/frontend.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
